# packages/jaraco.net/jaraco.net-7.1.0.tar.gz/jaraco.net-7.1.0/pinger/pinger.conf
[program:pinger]
command=/usr/bin/python -m jaraco.net.icmp 68.55.72.1
process_name=%(program_name)s
numprocs=1
directory=/home/pi
umask=022
priority=999
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
user=root
environment=PYTHONUSERBASE="/home/pi/.local"
