# packages/asus-charge-daemon/asus_charge_daemon-0.3.1-py3-none-any.whl/data/asuscharged.service
[Unit]
Description=ASUS Battery Charge Control Daemon

[Service]
Type=dbus
BusName=ca.cforrester.AsusChargeDaemon1
Restart=on-failure
RestartSec=1
User=root
StateDirectory=asuscharged
ConfigurationDirectory=asuscharged
Environment="PYTHONOPTIMIZE=2"
Environment="PYTHONUNBUFFERED=yes"
ExecStart=/usr/bin/asuscharged

[Install]
WantedBy=multi-user.target
